# Project-101/Third_Part/crontab
*/5 * * * * ./backup.sh   #every 5 minute

*/5 * * * * sudo /home/ec2-user/backup.sh